# codersBraggingRights/gunicorn.conf
description "Gunicorn application server handling codersBraggingRights"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu 
setgid www-data
chdir /home/ubuntu/codersBraggingRights

exec hackharvard/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/codersBraggingRights/codersBraggingRights.sock codersBraggingRights.wsgi:application
